# fw-zone-sync.service
[Unit]
Description=Sync Firewalld zone sources from DNS hostnames
After=network-online.target firewalld.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fw-zone-sync.sh -f /etc/fw-zone-sync/zones.ini
User=root
Group=root
Nice=10
IOSchedulingClass=idle
# Script writes JSON logs itself; keep journal for stdout/stderr
StandardOutput=journal
StandardError=journal
